# watchdog/com.sunrise.scanner-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sunrise.scanner-watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>-u</string>
    <string>watchdog/watchdog.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/tae/dev-mac-mini/launchpad-scanners</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/tae/dev-mac-mini/launchpad-scanners/watchdog/watchdog.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tae/dev-mac-mini/launchpad-scanners/watchdog/watchdog.log</string>
</dict>
</plist>
